# simul/cronjobs/crontab
*/10 * * * * /app/cronjobs/anomalies_e3.sh
*/25 * * * * /app/cronjobs/anomalies_e1.sh
*/57 * * * * /app/cronjobs/anomalies_e2.sh
